# luci-app-ipass/root/etc/uci-defaults/90_luci-ipass
#!/bin/sh
set -eu

ensure_site() {
  id="$1"
  name="$2"
  url="$3"
  category="$4"
  timeout="${5:-5}"

  if uci -q get "ipass.${id}" >/dev/null 2>&1; then
    return 0
  fi

  uci -q set "ipass.${id}=site"
  uci -q set "ipass.${id}.name=${name}"
  uci -q set "ipass.${id}.url=${url}"
  uci -q set "ipass.${id}.category=${category}"
  uci -q set "ipass.${id}.enabled=1"
  uci -q set "ipass.${id}.timeout=${timeout}"
}

if [ ! -s /etc/config/ipass ]; then
  cp /rom/etc/config/ipass /etc/config/ipass 2>/dev/null || true
fi

touch /etc/config/ipass

ensure_site baidu '百度' 'https://www.baidu.com/' '国内'
ensure_site taobao '淘宝' 'https://www.taobao.com/' '国内'
ensure_site weixin '微信' 'https://weixin.qq.com/' '国内'
ensure_site toutiao '今日头条' 'https://www.toutiao.com/' '国内'
ensure_site douyin '抖音' 'https://www.douyin.com/' '国内'
ensure_site xiaohongshu '小红书' 'https://www.xiaohongshu.com/' '国内'
ensure_site doubao '豆包' 'https://www.doubao.com/' '国内'
ensure_site xiaomi '小米' 'https://www.mi.com/' '国内'
ensure_site deepseek 'DeepSeek' 'https://www.deepseek.com/' '国内'
ensure_site kimi 'Kimi' 'https://kimi.moonshot.cn/' '国内'
ensure_site bilibili '哔哩哔哩' 'https://www.bilibili.com/' '国内'

if ! uci -q get ipass.google >/dev/null 2>&1 && ! uci -q get ipass.google_204 >/dev/null 2>&1; then
  ensure_site google 'Google' 'https://www.google.com/generate_204' '国际'
fi
ensure_site youtube 'YouTube' 'https://www.youtube.com/' '国际'
ensure_site apple 'Apple' 'https://www.apple.com/' '国际'
ensure_site tiktok 'TikTok' 'https://www.tiktok.com/' '国际'
ensure_site telegram 'Telegram' 'https://t.me/' '国际'
ensure_site chatgpt 'ChatGPT' 'https://chatgpt.com/' '国际'
ensure_site claude 'Claude' 'https://claude.ai/' '国际'
ensure_site gemini 'Gemini' 'https://gemini.google.com/' '国际'
ensure_site x 'X (Twitter)' 'https://x.com/' '国际'
ensure_site facebook 'Facebook' 'https://www.facebook.com/' '国际'
ensure_site github 'GitHub' 'https://github.com/' '国际'
ensure_site cloudflare 'Cloudflare' 'https://www.cloudflare.com/' '国际'

uci -q commit ipass
exit 0
